# usr/lib/systemd/user/filter-chain.service.d/10-biglinux-rt.conf
# BigLinux Microphone — give the upstream filter-chain.service the
# same realtime budget as our own biglinux-microphone-output.service.
#
# Without these limits the RT bump in
# `/usr/share/pipewire/filter-chain.conf.d/50-biglinux-priority.conf`
# can be silently capped by rtkit's per-process ceiling when the unit
# lacks `LimitRTPRIO` / `LimitNICE`. The values mirror the main
# `pipewire.service` budget shipped by the upstream PipeWire package,
# so behaviour stays consistent across all audio daemons on the box.
#
# `CapabilityBoundingSet` / `AmbientCapabilities` are deliberately
# absent: user-mode systemd cannot drop capabilities (the unit's
# spawn process is unprivileged), so configuring those keys produces
# `218/CAPABILITIES` and the service fails to start. RT scheduling
# here is granted by rtkit (CAP_SYS_NICE on rtkit-daemon), not by
# the unit itself, so dropping the cap directives is safe.
[Service]
LimitRTPRIO=95
LimitNICE=-19
LimitMEMLOCK=infinity
