# usr/lib/systemd/user/biglinux-mic-agc.service
[Unit]
Description=BigLinux Microphone AGC (automatic gain control)
After=pipewire.service wireplumber.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/biglinux-mic-agc
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
