# usr/lib/systemd/user/biglinux-microphone-echocancel.service
[Unit]
Description=BigLinux Microphone — acoustic echo cancellation (WebRTC AEC)
Documentation=https://github.com/biglinux/biglinux-noise-reduction-pipewire
PartOf=pipewire.service
After=pipewire.service wireplumber.service
Requisite=pipewire.service
Wants=pipewire.service wireplumber.service
# Only start when the user has actually opted into AEC — the conf is
# only written by the configurator when the toggle is on, and `pipewire
# -c` would crash-loop without it.
ConditionPathExists=%h/.config/pipewire/biglinux-microphone-echocancel.conf

[Service]
Type=simple
ExecStart=/usr/bin/pipewire -c %h/.config/pipewire/biglinux-microphone-echocancel.conf
Restart=on-failure
RestartSec=2s
RestartPreventExitStatus=254
StandardOutput=journal
StandardError=journal

# ── Realtime budget ──────────────────────────────────────────────────
# WebRTC AEC is sensitive to scheduling jitter — drops or late buffers
# manifest as residual echo. We mirror the rlimit budget the output
# unit uses so libpipewire-module-rt can honour the rt.prio it requests.
LimitRTPRIO=95
LimitNICE=-19
LimitMEMLOCK=infinity

# ── Hardening ────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectProc=invisible
PrivateTmp=true
RestrictNamespaces=true
RestrictRealtime=false
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=CAP_SYS_NICE
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_NETLINK
SystemCallArchitectures=native
SystemCallFilter=@system-service @resources
SystemCallErrorNumber=EPERM

[Install]
WantedBy=default.target
