# usr/lib/systemd/user/noise-reduction-pipewire.service
[Unit]
Description=PipeWire Input Filter Chain
After=pipewire.service
BindsTo=pipewire.service

[Service]
ExecStart=/usr/bin/pipewire-noise-remove
ExecStartPost=/usr/bin/pipewire-noise-remove start
ExecStopPost=/usr/bin/pipewire-noise-remove stop
Type=simple
Restart=on-failure

[Install]
WantedBy=pipewire.service
