# usr/lib/systemd/user/noise-reduction-pipewire.service
[Unit]
Description=Restore BigLinux Microphone Noise Reduction State
After=pipewire.service
Wants=pipewire.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'if [ -f "$HOME/.config/pipewire/filter-chain.conf.d/source-gtcrn-smart.conf" ]; then /usr/bin/pipewire-noise-remove start; fi'
RemainAfterExit=yes

[Install]
WantedBy=default.target
